# Cmder/config/profile.d/20-PSReadLine.ps1
Ensure-Module PSReadLine

Set-PSReadlineOption -ExtraPromptLineCount 1
Set-PSReadlineOption -EditMode Emacs -HistoryNoDuplicates
